# fan-control.service
[Unit]
After=network.target

[Service]
ExecStart=/usr/local/bin/control_fan

[Install]
WantedBy=default.target
